# website/freetron.service
[Unit]
Description=Freetron website daemon

[Service]
ExecStart=/usr/bin/freetron --daemon /srv/freetron

[Install]
WantedBy=multi-user.target
